# systemd/vaib-probe.timer
[Unit]
Description=Run vAIb-X lightweight probe every 2 minutes

[Timer]
OnBootSec=45s
OnUnitActiveSec=2min
Unit=vaib-probe.service
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
