# systemd/vaib-api.service
[Unit]
Description=vAIb-X backend API (port 4014)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/xsyprime/projects/vaiB-x
ExecStartPre=/home/xsyprime/projects/vaiB-x/scripts/vaib-prestart.sh
ExecStart=/usr/bin/env node server/api.mjs
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
