# systemd/vaib-probe.service
[Unit]
Description=vAIb-X lightweight probe + auto-restart guard
After=vaib-api.service vaib-ui.service

[Service]
Type=oneshot
WorkingDirectory=/home/xsyprime/projects/vaiB-x
ExecStart=/usr/bin/env bash /home/xsyprime/projects/vaiB-x/scripts/vaib-ensure-up.sh
